# docker/etc-supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10485760
logfile_backups=0
childlogdir=/var/log/supervisor
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
[program:hhvm1]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.1 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.1 --user root --mode daemon -vPidFile=/var/run/hhvm/pid1
[program:hhvm2]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.2 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.2 --user root --mode daemon -vPidFile=/var/run/hhvm/pid2
[program:hhvm3]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.3 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.3 --user root --mode daemon -vPidFile=/var/run/hhvm/pid3
[program:hhvm4]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.4 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.4 --user root --mode daemon -vPidFile=/var/run/hhvm/pid4
[program:hhvm5]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.5 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.5 --user root --mode daemon -vPidFile=/var/run/hhvm/pid5
[program:hhvm6]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.6 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.6 --user root --mode daemon -vPidFile=/var/run/hhvm/pid6
[program:hhvm7]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.7 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.7 --user root --mode daemon -vPidFile=/var/run/hhvm/pid7
[program:hhvm8]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.8 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.8 --user root --mode daemon -vPidFile=/var/run/hhvm/pid8
[program:hhvm9]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.9 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.9 --user root --mode daemon -vPidFile=/var/run/hhvm/pid9
[program:hhvm10]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.10 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.10 --user root --mode daemon -vPidFile=/var/run/hhvm/pid10
[program:hhvm11]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.11 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.11 --user root --mode daemon -vPidFile=/var/run/hhvm/pid11
[program:hhvm12]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.12 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.12 --user root --mode daemon -vPidFile=/var/run/hhvm/pid12
[program:hhvm13]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.13 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.13 --user root --mode daemon -vPidFile=/var/run/hhvm/pid13
[program:hhvm14]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.14 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.14 --user root --mode daemon -vPidFile=/var/run/hhvm/pid14
[program:hhvm15]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.15 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.15 --user root --mode daemon -vPidFile=/var/run/hhvm/pid15
[program:hhvm16]
directory=/opt/hhvm
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm.sock.16 -vRepo.Central.Path=/opt/hhvm/hhvm.hhbc.16 --user root --mode daemon -vPidFile=/var/run/hhvm/pid16
[program:nginx]
command=/usr/sbin/nginx
[program:nginx-agent]
command=/usr/bin/nginx-nr-agent.py -f -c /etc/nginx/agent.ini -p /tmp/agent.pid start
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
